# scripts/launchd/com.originpower.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.originpower.worker</string>
  <key>WorkingDirectory</key><string>__REPO__</string>
  <key>ProgramArguments</key><array>
    <string>__REPO__/.eval_venv/bin/python</string>
    <string>-m</string><string>studio</string><string>worker</string>
    
  </array>
  <key>EnvironmentVariables</key><dict>
    <key>STUDIO_DASH_TOKEN</key><string>__TOKEN__</string>
    <!-- launchd gives a minimal PATH; the pipeline shells out to Homebrew
         binaries (ffmpeg, blender, npx). Without this they are not found. -->
    <key>PATH</key><string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key><true/><key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>__REPO__/logs/worker.log</string>
  <key>StandardErrorPath</key><string>__REPO__/logs/worker.log</string>
</dict></plist>
